# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:webservice1]
command=python /app/webservice1/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/webservice1.err.log
stdout_logfile=/var/log/webservice1.out.log

[program:webservice2]
command=python /app/webservice2/app.py
autostart=true
autorestart=true
stderr_logfile=/var/log/webservice2.err.log
stdout_logfile=/var/log/webservice2.out.log
